# sh/profile.d/editor.sh
# Configure ViM to be the default editor
EDITOR=vim
export EDITOR
