# sh/profile.d/pager.sh
# Configure less to be the default pager
PAGER=less
export PAGER
